# roles/agg_pushgateway_exporter/templates/prom_aggregation_gateway.service.j2
[Unit]
Description=prom_aggregation_gateway
After=network-online.target

[Service]
Type=simple
ExecStart={{ prom_aggregation_gateway_binary }} --listen {{ prom_aggregation_gateway_listen_addr }}

[Install]
WantedBy=multi-user.target
