# roles/swift_exporter/templates/swift_exporter.service.j2
[Unit]
Description=Swift Exporter
After=network-online.target

[Service]
Type=simple

ExecStart={{ swift_exporter_binary }}

[Install]
WantedBy=multi-user.target

